# qubes-db/qubes-db-dom0.socket
[Socket]
ListenStream=/run/qubes/qubesdb.dom0.sock
Symlinks=/run/qubes/qubesdb.sock

[Install]
WantedBy=sockets.target
